# packages/base/busybox/fs-target/etc/init.d/rc.local
#!/bin/sh

# SEQ=98

. /etc/profile

case "$1" in
  start)
    if [ -x /etc/rc.local ] ; then
      /etc/rc.local
    fi
    ;;
esac

exit 0
